# deploy/night-shift-all.service
[Unit]
Description=Agent harness multi-product night_shift readiness (ORCH-P3)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=debian
Group=debian
# Absolute path avoids %h CHDIR failures seen 2026-07-17 (status=200/CHDIR)
WorkingDirectory=/home/debian/agent-harness
Environment=PRODUCT_VAULT_ROOT=/opt/second-brain/vault
Environment=WATCHLIST_VAULT_ROOT=/opt/second-brain/vault
Environment=NIGHT_SHIFT_PRODUCTS_FILE=/home/debian/agent-harness/config/night_shift_products.yaml
Environment=PATH=/home/debian/watchlist/.venv/bin:/home/debian/agent-harness/.venv/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/bin/bash -c 'PY="/home/debian/agent-harness/.venv/bin/python"; if [ ! -x "$PY" ]; then PY="/home/debian/watchlist/.venv/bin/python"; fi; if [ ! -x "$PY" ]; then PY=/usr/bin/python3; fi; exec "$PY" /home/debian/agent-harness/bin/night_shift_all_products.py'
SuccessExitStatus=0 1

[Install]
WantedBy=multi-user.target
