# deploy/night-shift-all.service
[Unit]
Description=Agent harness multi-product night_shift readiness (ORCH-P3)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=debian
Group=debian
# %h = home of User= (portable; no /home/<name> literals in unit text)
WorkingDirectory=%h/agent-harness
Environment=PRODUCT_VAULT_ROOT=/opt/second-brain/vault
Environment=WATCHLIST_VAULT_ROOT=/opt/second-brain/vault
Environment=NIGHT_SHIFT_PRODUCTS_FILE=%h/agent-harness/config/night_shift_products.yaml
# Prefer a product venv with pytest/mypy/ruff when present; else system PATH
Environment=PATH=%h/watchlist/.venv/bin:%h/agent-harness/.venv/bin:/usr/local/bin:/usr/bin:/bin
# Prefer harness venv python if installed; else system python3
ExecStart=/bin/bash -c 'PY="%h/agent-harness/.venv/bin/python"; if [ ! -x "$PY" ]; then PY="%h/watchlist/.venv/bin/python"; fi; if [ ! -x "$PY" ]; then PY=/usr/bin/python3; fi; exec "$PY" %h/agent-harness/bin/night_shift_all_products.py'
# Exit 1 = some product failed gates; still a completed overnight run
SuccessExitStatus=0 1

[Install]
WantedBy=multi-user.target
